# support/systemd/precached.service
[Unit]
Description=Process monitor and pre-caching daemon
Documentation=man:precached.service(8) man:precached.conf(5)

[Service]
ExecStart=/usr/sbin/precached -f -c /etc/precached/precached.conf
ExecReload=/bin/kill -HUP $MAINPID
Type=simple
StandardOutput=syslog
Restart=always

[Install]
WantedBy=multi-user.target
